# deploy/kallon-ptz-daemon.service.example
# Example systemd unit for Jetson / Ubuntu — copy to /etc/systemd/system/kallon-ptz-daemon.service
# and adjust user, paths, and secrets.
#
#   sudo cp deploy/kallon-ptz-daemon.service.example /etc/systemd/system/kallon-ptz-daemon.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kallon-ptz-daemon.service
#
# The values below match the bench Jetson layout (user `khalifa`, repo at
# /home/khalifa/kallon). Edit if you clone the repo elsewhere or use a
# dedicated service account.

[Unit]
Description=Kallon ONVIF PTZ daemon (JSON/TCP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=khalifa
Group=khalifa
WorkingDirectory=/home/khalifa/kallon

# Required on production systems (overrides script default password).
Environment=CAMERA_PASSWORD=change_me

# Listen only on loopback so the daemon is reachable from ArtemisOS or other
# local clients on the same host, but not from the network.
Environment=PTZ_LISTEN_HOST=127.0.0.1
Environment=PTZ_LISTEN_PORT=8765

# Camera ONVIF connection — edit IP/port/user as needed for this device.
ExecStart=/usr/bin/python3 /home/khalifa/kallon/kallon_ptz_daemon.py \
    --host 192.168.1.108 -P 80 -u admin \
    --listen-host 127.0.0.1 --listen-port 8765

Restart=on-failure
RestartSec=3

# Hardening (tune for your image).
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
