# deploy/kallon-ptz-daemon.service.example
# Example systemd unit for Jetson / Ubuntu — copy to /etc/systemd/system/kallon-ptz-daemon.service.
# The installer (scripts/install/80-watchdogs.sh) renders the real unit;
# this file is a reference for manual installs or inspection only.
#
# Replace RUNTIME_USER with your actual runtime user (default: khalifa).
#
#   sudo cp deploy/kallon-ptz-daemon.service.example /etc/systemd/system/kallon-ptz-daemon.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kallon-ptz-daemon.service

[Unit]
Description=Kallon ONVIF PTZ daemon (JSON/TCP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=RUNTIME_USER
Group=RUNTIME_USER
WorkingDirectory=/opt/kallon

# Cameras are read from device.env (CAMERA_IPS, CAMERA_RTSP_USER,
# CAMERA_PASSWORD, CAMERA_ONVIF_PORT) and exposed 1-based as camera 1..N,
# aligned with mediamtx paths cam1..camN. This is the production path; do not
# hard-code a single --host here. The installer renders the real unit from
# scripts/install/80-watchdogs.sh.
EnvironmentFile=/etc/kallon/device.env

# Listen only on loopback so the daemon is reachable from ArtemisOS, the tower
# lab dashboard, or other local clients on the same host — never the network.
Environment=PTZ_LISTEN_HOST=127.0.0.1
Environment=PTZ_LISTEN_PORT=8765

ExecStart=/usr/bin/python3 /opt/kallon/kallon_ptz_daemon.py \
    --listen-host 127.0.0.1 --listen-port 8765

Restart=on-failure
RestartSec=3

# Hardening (tune for your image).
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
