# deploy/kallon-watchdog.service.example
# Example systemd unit for the Kallon health & tamper watchdog.
# The installer (scripts/install/80-watchdogs.sh) renders the real unit; this
# file is a reference for manual installs or inspection only.
#
# Replace RUNTIME_USER with your Jetson login (set RUNTIME_USER in /etc/kallon/device.env).
#
#   sudo cp deploy/kallon-watchdog.service.example /etc/systemd/system/kallon-watchdog.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kallon-watchdog.service

[Unit]
Description=Kallon health and tamper watchdog
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
User=RUNTIME_USER
Group=RUNTIME_USER
# The user must also be in the `gpio` and `i2c` groups for /dev/gpiochip* and /dev/i2c-*.
SupplementaryGroups=gpio i2c

WorkingDirectory=/opt/kallon

# Per-device secrets and configuration. Written by the installer.
EnvironmentFile=/etc/kallon/device.env

# Jetson.GPIO needs to know the board model. The installer auto-detects this
# from /proc/device-tree/model; override here only if auto-detection fails.
# Environment=JETSON_MODEL_NAME=JETSON_ORIN_NANO

ExecStart=/usr/bin/python3 /opt/kallon/kallon_watchdog.py

Restart=on-failure
RestartSec=3

# Hardening (loosen if you hit a path the daemon legitimately needs).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadOnlyPaths=/etc/kallon
# /dev access for I2C and GPIO must remain available; do not set PrivateDevices=true.

[Install]
WantedBy=multi-user.target
